# pass/.config/systemd/user/pass-secret-service.service
[Unit]
Description=Pass Secret Service

[Service]
ExecStart=/usr/bin/pass-secret-service -n
Restart=on-failure

[Install]
WantedBy=default.target
